# rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
sh -c "source /root/.profile && python /root/sources/Omega2Plus/test.py" &
sleep 1
disown
exit 0
